# scripts/crontab
* * * * * python /scripts/src/download_ticker_history.py COINBASE BTC-USD
